# plex-grouch.timer
[Unit]
Description=Oscar the Plex Grouch - Hourly Cleanup

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
